# config/openvox-repo-sync.service
[Unit]
# Mirrors the OpenVox / Vox Pupuli package repositories into
# /opt/openvox-pkgs so the openvox-gui server can act as a local
# package mirror. Triggered by openvox-repo-sync.timer (nightly) or
# manually via `systemctl start openvox-repo-sync.service` or the
# Installer page in the GUI.
Description=OpenVox repository mirror sync
Documentation=file:///opt/openvox-gui/docs/INSTALLER.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Run as root so curl can write into /opt/openvox-pkgs and chown the
# resulting tree to the puppet user. The script itself does not need
# privileged operations beyond filesystem writes.
User=root
Group=root

# Path to the sync script. Installed by openvox-gui's install.sh into
# the install directory's scripts/ subfolder.
ExecStart=/opt/openvox-gui/scripts/sync-openvox-repo.sh --quiet

# Reasonable timeout: full mirror (yum + apt + windows + mac, both
# v7 and v8, both arches) can take a while on first run -- empirically
# the yum portion alone runs ~1h40m, apt is the largest tree, and the
# old 2h cap killed in-progress full syncs at the apt stage (verified
# 2026-04-23 on test server: yum finished at 17:01Z having started at
# 15:20Z, then SIGTERM at 17:20Z mid-apt). 6h gives comfortable
# headroom for a cold full sync; incremental runs finish in minutes.
TimeoutStartSec=6h

# Don't restart on failure -- timer will retry on next schedule
Restart=no

# Logging
StandardOutput=journal
StandardError=journal

# Environment overrides can be placed in /etc/sysconfig/openvox-repo-sync
# (RHEL family) or /etc/default/openvox-repo-sync (Debian family).
# The openvox-gui .env is loaded so proxy settings configured via the
# GUI's Proxy Configuration page are available to the sync script.
EnvironmentFile=-/etc/sysconfig/openvox-repo-sync
EnvironmentFile=-/etc/default/openvox-repo-sync
EnvironmentFile=-/opt/openvox-gui/config/.env

[Install]
WantedBy=multi-user.target
