# config/openvox-gui.service
[Unit]
Description=OpenVox GUI - Puppet Management Web Interface
After=network.target puppetserver.service puppetdb.service
Wants=puppetdb.service

[Service]
Type=simple
User=puppet
Group=puppet
WorkingDirectory=INSTALL_DIR/backend
EnvironmentFile=INSTALL_DIR/config/.env
ExecStart=INSTALL_DIR/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 4567
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
PrivateTmp=false
ReadWritePaths=INSTALL_DIR/data INSTALL_DIR/logs /opt/puppetlabs/puppet/cache /etc/puppetlabs/code/environments /tmp

[Install]
WantedBy=multi-user.target
